# vps-gateway/deploy/systemd/sandalphone-vps-gateway.service
[Unit]
Description=Sandalphone VPS Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sandalphone/vps-gateway
EnvironmentFile=/opt/sandalphone/vps-gateway/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=2
User=sandalphone
Group=sandalphone

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/sandalphone/vps-gateway

[Install]
WantedBy=multi-user.target
